# src/production/mariadb/fedora/systemd/mariadb-production-fedora-maintenance.timer
[Unit]
Description=Perform routine mariadb maintenance

[Timer]
OnStartupSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
